# init/nis.conf
# compatibility job for sysvinit dependencies

start on starting ypserv
